# docs/ubuntu_desktop/bootstrap.service
[Unit]
Description=Run Bootstrap Script

[Service]
ExecStart=/usr/local/bin/first_boot.sh
User=root
Type=oneshot
RemainAfterExit=no

[Install]
WantedBy=multi-user.target