# docs/ubuntu_daita/bootstrap.service
[Unit]
Description=Run Bootstrap Script

[Service]
ExecStart=/usr/local/bin/first_boot.sh
User=root
Type=oneshot
RemainAfterExit=true

[Install]
WantedBy=multi-user.target